# Docker/PyCrontab/PyCron-test-db/crontab
#run python script every minutes
*/6 * * * * /usr/local/bin/python /app/TGtoMariadb_All.py > /proc/1/fd/1 2>/proc/1/fd/2
